# conf/app.socket
[Unit]
Description=app socket

[Socket]
ListenStream=/run/app.sock

[Install]
WantedBy=sockets.target